# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini /usr/local/app/routeapi/routeapi.ini

[program:nginx]
command=nginx -g "daemon off;"
